# .docker/bin/webserver/crontab
* * * * * root echo helloxxx > /dev/pts/0
* * * * * root date > /dev/pts/0
* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1 

